# templates/complex.service.j2
[Unit]
Description=Complex Application

[Service]
User=root
Group=root

ExecStart=/opt/complex/main.py hello
WorkingDirectory=/opt/complex/

Restart=always

[Install]
WantedBy=multi-user.target
